# retired-launchagents/com.crewclock.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.crewclock.tunnel</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/maximilian/Desktop/crew-clock/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--url</string>
    <string>http://localhost:3000</string>
    <string>--no-autoupdate</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/maximilian/Desktop/crew-clock/logs/tunnel.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/maximilian/Desktop/crew-clock/logs/tunnel.log</string>
</dict>
</plist>
